# demarkus.service
[Unit]
Description=demarkus server
After=network.target

[Service]
Type=simple
User=demarkus
ExecStart=/srv/demarkus/demarkus-server
Environment=DEMARKUS_ROOT=/srv/demarkus/blog
Environment=DEMARKUS_TLS_CERT=/srv/demarkus/certs/fullchain.pem
Environment=DEMARKUS_TLS_KEY=/srv/demarkus/certs/privkey.pem
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
